# twitter/cron/crontab
MAILTO=twitter
@hourly /opt/app/twitter/bin/run